# systemd/comunito-portal-v8.service
[Unit]
Description=Comunito Portal V8
After=network-online.target
Wants=network-online.target

[Service]
User=__USER__
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/comunito-venv/bin/python __APP_DIR__/app/app/portal_v8_main.py
Restart=always
RestartSec=3
Environment=TZ=America/Mexico_City
Environment=OMP_NUM_THREADS=2
Environment=OPENBLAS_NUM_THREADS=1
Environment=PYTHONPATH=__APP_DIR__

[Install]
WantedBy=multi-user.target
